# agent/ptpbox-agent.service
[Unit]
Description=PTPBox observation API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@PTPBOX_USER@
Group=@PTPBOX_GROUP@
SupplementaryGroups=clock
Environment=PTPBOX_ROOT=@PTPBOX_ROOT@
Environment=PTPBOX_STATE_DIR=@PTPBOX_ROOT@/runtime
Environment=PTPBOX_LOG_DIR=/var/log/ptpbox
Environment=PTPBOX_TOPOLOGY=/etc/ptpbox/topology.json
Environment=PTPBOX_BIND=0.0.0.0
Environment=PTPBOX_PORT=8090
Environment=PTPBOX_WEB_ROOT=/opt/ptpbox-web/static
ExecStart=/usr/bin/python3 /opt/ptpbox-web/agent/ptpbox_agent.py
Restart=on-failure
RestartSec=2
KillMode=process
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=@PTPBOX_ROOT@/runtime /etc/linuxptp /run/netns /run/ptpbox /var/log/ptpbox

[Install]
WantedBy=multi-user.target
